# scripts/systemd/vyrdon-lab-team-room.timer
[Unit]
Description=VYRDON Lab Console - Daily Team Room update

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target

